# config/autostart/fcitx5.desktop
[Desktop Entry]
Exec=/usr/bin/fcitx5
Name=fcitx5
Type=Application
