# config/autostart/greenclip.desktop
[Desktop Entry]
Exec=/usr/bin/greenclip daemon
Name=greenclip
Type=Application
