# config/autostart/nitrogen.desktop
[Desktop Entry]
Exec=/usr/bin/nitrogen --restore
Name=nitrogen
Comment=X11 background previewer and setter
Type=Application
