# config/autostart/nm-applet.desktop
[Desktop Entry]
Exec=/usr/bin/nm-applet
Name=nm-applet
Comment=network monitor and control GUI applet
Type=Application
