# config/autostart/polkit-dumb-agent.desktop
[Desktop Entry]
Exec=/usr/bin/polkit-dumb-agent
Name=polkit-dumb-agent
Type=Application
