# config/autostart/poweralertd.desktop
[Desktop Entry]
Exec=/usr/bin/poweralertd
Name=poweralertd
Comment=power alert notifications
Type=Application
