# config/autostart/volumeicon.desktop
[Desktop Entry]
Exec=/usr/bin/volumeicon
Name=volumeicon
Type=Application
