# config/autostart/xcompmgr.desktop
[Desktop Entry]
Exec=/usr/bin/xcompmgr -c -C -t-5 -l-5 -r4.2 -o.55
Name=xcompmgr
Comment=sample X compositing manager
Type=Application
